# systemd_user/notify_addr.timer
[Unit]
Description=v6 logging timer

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m
Unit=notify_addr.service

[Install]
WantedBy=timers.target
